# celery_worker.ini
[program:celery_worker]
command=/data/xsadmin_project/env/bin/celery -A xsadmin worker -l info
directory=/data/xsadmin_project/xsadmin
user=root
numprocs=1
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/data/xsadmin_project/logs/celery_worker.log
