# conf/supervisor_uwsgi_xsadmin.ini
[program:uwsgi_xsadmin]
command=/data/xsadmin_project/env/bin/uwsgi --ini uwsgi.ini
directory=/data/xsadmin_project/xsadmin
startsecs=0
stopwaitsecs=0
autostart=true
autorestart=true
user=root
redirect_stderr=true
stdout_logfile = /data/xsadmin_project/logs/xsadmin.log