# conf/supervisor_celery_worker.ini
[program:celery_worker]
command=/data/xsadmin_deploy/env/bin/celery -A xsadmin worker -l info
directory=/data/xsadmin_deploy/xsadmin
user=root
numprocs=1
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/data/xsadmin_deploy/logs/celery_worker.log
