# conf/supervisor_celery_beat.ini
[program:celery_beat]
command=/data/xsadmin_project/env/bin/celery -A xsadmin beat -l info
directory=/data/xsadmin_project/xsadmin
user=root
numprocs=1
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/data/xsadmin_project/logs/celery_beat.log
